# smart-trade.service
[Unit]
Description=Claude Smart Trade Bot v3.0
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/smart-trade
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/smart-trade/venv/bin/python /opt/smart-trade/smart_trade_v3.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/smart-trade/bot.log
StandardError=append:/var/log/smart-trade/error.log

[Install]
WantedBy=multi-user.target
